# tags/kamikaze_7.09/package/iw-olsr/files/etc/crontabs/root
* * * * * /usr/lib/webif/LUA/olsr/minute.cron


