# package/webif-iw-lua-chillispot/files/usr/lua/lua-xwrt/pkgs/chilli/minute.cron
#!/usr/bin/lua
package.cpath = "?;./?.so;/usr/share/internet-wifi/lib.so/?.so;" .. package.cpath 
package.path = "?;./?.lua;/usr/share/internet-wifi/?.lua;/usr/share/internet-wifi/lib/?.lua;/usr/share/internet-wifi/pkgs/?.lua;" .. package.path 
-- Common Functions
require("common")

function must_be_running()
  if io.exists("/var/run/chilli.pid") == false then
    os.execute("/usr/sbin/chilli")
  end
end

must_be_running()
