# tags/kamikaze_7.09/package/webif-lua-olsr/files/usr/lib/lua/lua-xwrt/pkgs/olsr/minute.cron
#!/usr/bin/lua
package.path = package.path .. ";/usr/lib/webif/LUA/?.lua;/usr/share//lua/5.1/?.lua"
-- Common Functions
require("common")
local output_file = "/etc/olsrd.conf"
local previus load_fiel(output_file)
local prev_len = string.len(previus)
dofile("/usr/lib/webif/LUA/olsr/parsers")
local newfile load_file(output_file)
local new_len = string.len(newfile)
if prev_len ~= new_len then 
  os.execute("/etc/init.d/olsrd restart")
end  
 