# whiterussian/package/miniupnpd/files/etc/hotplug.d/iface/10-miniupnpd
#!/bin/sh
[ "$ACTION" = "ifup" ] && [ "$INTERFACE" = "wan" ] && /etc/init.d/S??miniupnpd start
