# whiterussian/package/webif/files/etc/hotplug.d/iface/10-opendns
#!/bin/sh
[ "$INTERFACE" = "wan" -a -x /etc/init.d/S??opendns ] && {
	case "$ACTION" in
		ifup)	/etc/init.d/S??opendns enable ;;
		ifdown)	/etc/init.d/S??opendns disable ;;
	esac
}
